# scripts/launchd/com.nathanieldavis.paper-trading.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nathanieldavis.paper-trading</string>

    <!--
        Pre-stage data that doesn't depend on the midnight-UTC Binance close
        (macros, CoinMetrics, HL instruments, Binance OI/LSR, premium index,
        volume tail). Runs daily ahead of the automated 00:05 UTC full
        advisory run (com.nathanieldavis.paper-trading-full) so the full run
        can pass --skip-prestage and avoid duplicate work.

        Previously ran `daily_paper_run.py --non-binance-only --notify`,
        which only touched 3 of the 6 prestage outputs — `--non-binance-only`
        explicitly skips Binance-derived steps (OI/LSR, volume, premium) at
        `_oi_lsr_step` etc. in daily_paper_run.py. As a result OI never got
        refreshed by the aux flow, and the full run's freshness check
        correctly fell back to a full prestage on every fire. Switched to
        `prestage_daily.py` — which is purpose-built for this — so all 6
        outputs are touched and --skip-prestage actually fires.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/python@3.10/bin/python3.10</string>
        <string>/Users/nathanieldavis/pysystemtrade-prod/scripts/prestage_daily.py</string>
        <string>--config</string>
        <string>/Users/nathanieldavis/pysystemtrade-prod/config/crypto_perps_1k.yaml</string>
        <string>--env</string>
        <string>prod</string>
    </array>

    <!--
        TZ=UTC is exported into the running script so date strings in logs are UTC,
        but it does NOT affect when launchd schedules the job — see comment on
        StartCalendarInterval below.
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>UTC</string>
    </dict>

    <!--
        macOS launchd evaluates StartCalendarInterval in SYSTEM LOCAL TIME, not
        the TZ env var. System tz is America/New_York, so:
            Hour=18, Minute=30 → fires at 18:30 EDT (= 22:30 UTC) during DST,
                                  and 18:30 EST (= 23:30 UTC) during winter.

        We accept the 1-hour drift between EDT/EST. Moved from 15:00 to 18:30
        on 2026-05-18 so the aux prestage lands AFTER US-market settlements
        each day. yfinance-backed feeds (macro_factors: SPX/DXY/US10Y/US5Y/
        gold/VIX/oil; etf_flows: IBIT/ETHA) only finalize once the underlying
        cash/futures sessions close at ~21:00 UTC EDT / ~22:00 UTC EST.
        Firing at 15:00 ET pulled partial intraday bars or stale T-1 closes
        and the main 00:05 UTC full run skipped re-pull (aux <12h old), so
        the live system traded on a macro snapshot up to 1 trading day behind
        the latest available US close. 18:30 ET buys 1.5h (DST) / 0.5h (EST)
        of buffer ahead of the 00:05 UTC full-run fire — tight in EST winter
        if an aux outlier runs >30 min, but typical aux runtime is 5-15 min.
    -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/nathanieldavis/pysystemtrade-prod/live/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nathanieldavis/pysystemtrade-prod/live/launchd_stderr.log</string>

    <!-- Do not fire immediately on load — wait for next scheduled time. -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Do not restart if it exits — daily one-shot. -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
